# supervisord.conf
[program:node]
command = NODE_ENV=production node app.js
directory = /home/dotcloud/current
