# supervisord.conf
[program:node]
command = export NODE_ENV='production';node app.js
directory = /home/dotcloud/current
