# supervisord.conf
[program:node]
command = sh -c 'NODE_ENV=production node app.js'
directory = /home/dotcloud/current
